# tools/release_templates/systemd/assistant-task-api.service.template
[Unit]
Description=HarborBeacon Assistant Task API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__INSTALL_ROOT__/current
EnvironmentFile=-__ENV_FILE__
ExecStart=__INSTALL_ROOT__/current/templates/bin/run-assistant-task-api
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
